# files/nginx.timer
[Unit]
Description=Delays nginx to avoid race condition with DNS

[Timer]
OnBootSec=1min
Unit=nginx.service

[Install]
WantedBy=basic.target
